# systemd/stackchan-voice-loop.service
[Unit]
Description=StackChan wake-word voice loop (say "Hey Cherub" -> prompt routes to OpenClaw main session)
Documentation=file:///home/sunkencity999/.openclaw/workspace/stackchan/README.md
After=stackchan-gateway.service network-online.target
Wants=stackchan-gateway.service
PartOf=stackchan-gateway.service
# Mutually exclusive with the phase-B convo loop: both poll the same mic.
Conflicts=stackchan-convo.service
OnFailure=unit-failure-alert@%N.service

[Service]
Type=simple
WorkingDirectory=/home/sunkencity999/.openclaw/workspace
ExecStart=/home/sunkencity999/.local/share/uv/tools/stackchan-mcp/bin/python /home/sunkencity999/.openclaw/workspace/scripts/stackchan_voice_loop.py run
Restart=on-failure
RestartSec=15s
Nice=5
CPUWeight=60
StandardOutput=append:/home/sunkencity999/.openclaw/workspace/stackchan/voice_loop/voice.log
StandardError=append:/home/sunkencity999/.openclaw/workspace/stackchan/voice_loop/voice.log

[Install]
WantedBy=default.target
